# barber-bot.service
[Unit]
Description=Barber Reminder Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/barber
ExecStart=/opt/barber/venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
